# postfix_supervisord.conf
[program:postfix]
directory   = /etc/postfix
command     = /usr/local/bin/postfix_proxy.sh
startsecs   = 0
autorestart = false


[program:rsyslog]
command=/usr/sbin/rsyslogd -n 
